# avbot.service.in
[Unit]
Description=avbot bridges QQ/IRC/XMPP world together
Wants=network-online.target

[Service]
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/avbot --nostdin
RestartSec=1s
Type=simple
Environment=LANG=en_US.UTF-8
Environment=TZ=Asia/Shanghai

[Install]
WantedBy=multi-user.target
